# resources/etc/service/sv_mysql/run
#!/bin/bash

su mysql -c '/usr/sbin/mysqld' -s '/bin/bash' &
